# supervisord.ini
[program:InterestBlog]
command=python -O /home/taochengwei/git/interest.blog/Product.py
autostart=true
autorestart=true
startsecs=5
priority=1
stopasgroup=true
killasgroup=true
